# systemd-services/nerdbot-flask.service
[Unit]
Description=NerdBot Flask Backend Server
Documentation=file:///home/mark/nerdbot-backend/flask_server/server.py
After=network.target sound.target
Wants=nerdbot-fan.service

[Service]
Type=simple
User=mark
Group=mark
WorkingDirectory=/home/mark/nerdbot-backend
Environment=PATH=/home/mark/.local/bin:/home/mark/.nvm/versions/node/v20.17.0/bin:/home/mark/.cargo/bin:/home/mark/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/local/games:/usr/games
ExecStartPre=/bin/bash -c 'if aplay -l | grep -q Anker; then CARD_NUM=$(aplay -l | grep Anker | sed -n "s/card \\([0-9]\\):.*/\\1/p"); if [ -n "$CARD_NUM" ]; then amixer -c $CARD_NUM cset numid=3 95 2>/dev/null; fi; fi; amixer sset Master 75%% 2>/dev/null || amixer sset PCM 75%% 2>/dev/null'
ExecStart=/bin/bash -c 'source setup_env.sh && python -m flask_server.server'
Restart=always
RestartSec=10
StandardOutput=append:/home/mark/nerdbot.log
StandardError=append:/home/mark/nerdbot.log

# Allow access to audio devices and GPIO
SupplementaryGroups=audio video gpio spi i2c input render

[Install]
WantedBy=multi-user.target